# src/supervisor.d/queue.ini
[program:swoole_queue]
command=/usr/bin/php /datas/www/testdev/artisan wangpkg:swoole_queue xMaster
user=root
autostart=true
autorestart=true
startretries=10
stderr_logfile=/datas/logs/supervisor.err.swoole_queue.log
stdout_logfile=/datas/logs/supervisor.success.swoole_queue.log

[program:normal_queue]
command=/usr/bin/php /datas/www/testdev/artisan wangpkg:swoole_queue xSlave test 10
user=root
autostart=true
autorestart=true
startretries=10
stderr_logfile=/datas/logs/supervisor.err.normal_queue.log
stdout_logfile=/datas/logs/supervisor.success.normal_queue.log
